# files/usr/local/etc/xdg/autostart/nss-trust-root-ca.desktop.desktop
[Desktop Entry]
Type=Application
Name=Add site root CA to user NSS database.
Exec=/usr/local/libexec/nss-trust-root-ca
StartupNotify=false
NoDisplay=true
