# files/etc/cron.d/davical.dav_server
MAILTO=root
0,15,30,45 * * * * ${nginx_user} ${davical_repo_dir}/scripts/cron-sync-ldap.php ${fqdn}
