# files/etc/cron.d/freeradius.radius_server
MAILTO=root
@daily ${freeradius_user} find ${freeradius_tlscache_dir} -mindepth 1 -mtime +2 -exec rm -f {} +
