# files/etc/cron.d/invidious.invidious_server
MAILTO=root
0  3 * * * root /usr/local/libexec/invidious-update           -q ${invidious_local_username} ${invidious_repo_dir}
30 3 * * * root /usr/local/libexec/invidious-companion-update -q ${invidious_local_username} ${invidious_companion_repo_dir}
0  4 * * * root service invidious status > /dev/null && service invidious restart > /dev/null
