# files/etc/cron.d/zfs-trim.freebsd
@weekly root zpool list -Ho name | xargs -r -n1 zpool trim
